# xos/synchronizer/supervisor/monitoring_observer.conf
[program:monitoring_service_observer]
command=python /opt/xos/synchronizers/monitoring/monitoring_synchronizer.py -C /opt/xos/synchronizers/monitoring/monitoring_synchronizer_config
